# Extracted Makerbot Mini Firmware/ubifs-root/etc/profile.d/librt_hack.sh
if [ "$PS1" ]; then
    export LD_PRELOAD=librt.so.1
fi
